# recipes-kratos/kratos/files/kratos.service
[Unit]
Description=Kratos service
After=can.service

[Service]
ExecStart=/usr/lib/kratos/kratos run

[Install]
WantedBy=multi-user.target
